# telegraf/rootfs/etc/services.d/telegraf/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Telegraf
# Runs the Telegraf Server
# ==============================================================================

bashio::log.info 'Starting Telegraf...'

# Run Telegraf
exec telegraf
