# telegraf/rootfs/etc/services.d/telegraf/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: Telegraf
# Runs the Telegraf Server
# ==============================================================================

bashio::log.info 'Starting Telegraf...'

# Run Telegraf
# exec telegraf

set -e

if [ "${1:0:1}" = '-' ]; then
    set -- telegraf "$@" --config "/etc/telegraf/telegraf.conf" --config-directory "/config/telegraf"
fi

if [ "$(id -u)" -ne 0 ]; then
    exec "$@" --config "/etc/telegraf/telegraf.conf" --config-directory "/config/telegraf"
else
    # Allow telegraf to send ICMP packets and bind to privliged ports
    setcap cap_net_raw,cap_net_bind_service+ep /usr/bin/telegraf || echo "Failed to set additional capabilities on /usr/bin/telegraf"

    exec su-exec telegraf "$@" --config "/etc/telegraf/telegraf.conf" --config-directory "/config/telegraf"
fi
